# test-gui-04/root/etc/skel/.config/autostart/sforzato.desktop
[Desktop Entry]
Type=Application
Name=setup sfor
Path=/usr/share/sforzato
Exec=setup_sforzato.sh
Icon=
Terminal=false
